# examples/lxc-backup@.timer
[Unit]
Description=Daily timer to run LXC backup for %i

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
